# supervisord.conf
[supervisord]
nodaemon=true

[program:infinidb]
command=./infinidb_non_demonized_healper.sh
