# supervisord.conf
[supervisord]
nodaemon=true

[program:infinidb]
command=service infinidb start
